# files/xvfb.conf
[program:xvfb]
command=Xvfb :1 -screen 0 1024x768x16
user=android
autostart=true
autorestart=true
